# target_files/obmc-bmc-quiesce.target
[Unit]
Description=BMC Quiesce Target
Wants=multi-user.target
After=multi-user.target
