# target_files/obmc-drive-powered-off@.target
[Unit]
Description=Drive%i (Powered Off)
Conflicts=obmc-drive-powered-on@%i.target
