# target_files/obmc-drive-hard-reboot@.target
[Unit]
Description=Hard Reboot Drive%i
After=multi-user.target
Wants=obmc-drive-powercycle@%i.target
After=obmc-drive-powercycle@%i.target
