# target_files/obmc-bmc-active.target
[Unit]
Description=Active BMC Target
BindsTo=phosphor-set-bmc-active-passive.service
