# target_files/obmc-fan-control-ready@.target
[Unit]
Description=Fan Control Ready
Conflicts=obmc-chassis-powered-off@%i.target
