# target_files/obmc-chassis-powercycle@.target
[Unit]
Description=Chassis%i Cycle
After=multi-user.target

