# target_files/obmc-host-crash@.target
[Unit]
Description=Host instance %i crashed
Conflicts=obmc-chassis-poweroff@%i.target
After=obmc-chassis-poweroff@%i.target
