# target_files/obmc-drive-powered-on@.target
[Unit]
Description=Drive%i (Powered On)
Conflicts=obmc-drive-powered-off@%i.target
