# target_files/obmc-chassis-blackout@.target
[Unit]
Description=Chassis Blackout Target
After=multi-user.target
Conflicts=obmc-chassis-poweron@.target
