# target_files/obmc-drive-powercycle@.target
[Unit]
Description=Drive%i (Power Cycle)
After=multi-user.target
