# target_files/obmc-bmc-service-failure.target
[Unit]
Description=BMC service entered failure state
RefuseManualStop=yes
