# target_files/obmc-host-shutdown@.target
[Unit]
Description=Power%i Host Off
After=multi-user.target
