# target_files/obmc-host-shutdown@.target
[Unit]
Description=Power%i Host Off
After=multi-user.target
OnFailure=obmc-chassis-poweroff@%i.target
OnFailureJobMode=flush
