# target_files/obmc-fan-control.target
[Unit]
Description=Fan Control
RefuseManualStart=yes
