# target_files/obmc-host-hardpoweroff@.target
[Unit]
Description=Power%i Host Hard Power Off
After=multi-user.target
