# target_files/obmc-bmc-service-quiesce@.target
[Unit]
Description=BMC Quiesce Target
RefuseManualStop=yes
