# target_files/obmc-host-release-pre@.target
[Unit]
Description=The Host%i is ready to be prepared for release
After=phosphor-ipmi-host.service
RefuseManualStop=yes
