# target_files/obmc-host-graceful-shutdown@.target
[Unit]
Description=Power%i Host Graceful Shutdown
After=multi-user.target
OnFailure=obmc-chassis-poweroff@%i.target
OnFailureJobMode=flush
