# target_files/obmc-host-graceful-quiesce@.target
[Unit]
Description=Graceful Quiesce Target
After=multi-user.target
Conflicts=obmc-chassis-poweroff@%i.target
