# target_files/obmc-host-reset@.target
[Unit]
Description=Host%i (Reset Check)
Conflicts=obmc-host-stop@%i.target

[Install]
WantedBy=multi-user.target
