# service_files/phosphor-chassis-wait-for-smp-poweron.service
[Unit]
Description=Wait for all SMP chassis to power on
# Block the chassis 0 power-on target from completing
Before=obmc-power-on@0.target
Wants=obmc-power-on@0.target
# Ensure chassis 0 state manager is running (it monitors the others)
After=xyz.openbmc_project.State.Chassis0.service
Wants=xyz.openbmc_project.State.Chassis0.service
Conflicts=obmc-chassis-poweroff@0.target

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/usr/bin/phosphor-chassis-wait-for-smp-poweron
TimeoutStartSec=300
SyslogIdentifier=chassis-wait-smp-poweron

[Install]
RequiredBy=obmc-chassis-poweron@0.target