# service_files/phosphor-bmc-quiesce-reboot.service
[Unit]
Description=Reboot BMC after entering BMC Quiesced
DefaultDependencies=no

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl reboot
